# modules/xbmc/files/standalone/dotfiles/.config/autostart/dconf-settings.desktop
[Desktop Entry]
Name=Dconf Settings
Comment=Fixes shit
Exec=~/.config/dconf-settings.sh
Terminal=false
Type=Application
