# modules/monitoring/templates/grafana.service.erb
[Unit]
Description=Grafana service
After=network.target

[Service]
User=grafana
Group=grafana
WorkingDirectory=/usr/share/grafana
ExecStart=/usr/sbin/grafana-server --config=<%= @grafana_config %> cfg:default.paths.data=/var/lib/grafana cfg:default.paths.logs=/var/log/grafana
LimitNOFILE=10000
TimeoutStopSec=20
SuccessExitStatus=0 2

[Install]
WantedBy=multi-user.target
