# modules/owncloud/files/backup-owncloud-db.cron
#!/bin/bash
BACKUPDIR=/usr/share/webapps/owncloud/database-backups

mysqldump owncloud > "$BACKUPDIR/$(date +%Y-%m-%d)-owncloud.sql"

find $BACKUPDIR/* -mtime +14 -maxdepth 0 -delete
