# modules/puppetmaster/files/puppet.cron
*/5 * * * * /usr/local/bin/update_puppet
