# modules/puppetmaster/files/puppet.cron
*/5 * * * * root /usr/local/bin/update_puppet
